# openc2x/files/etc/uci-defaults/uci-defaults-openc2x
#!/bin/sh

/etc/init.d/openc2x enabled && {
    /etc/init.d/openc2x stop
    /etc/init.d/openc2x disable
}

exit 0 
